# ansible/roles/supervisor/templates/osrm.conf
[program:osrm]
directory=/home/osrm
command=/usr/local/bin/osrm-routed -p 5000 berlin-latest.osrm
user=osrm
