# ansible/roles/finkcrow/templates/finkcrow_supervisor.conf
[program:finkcrow]
command=/bin/bash /opt/finkcrow.sh
